# steam.target
[Unit]
Description=Steam Headless Shim
Requires=steam.service

[Install]
WantedBy=default.target